# bin/Release/netcoreapp2.2/ubuntu.18.10-x64/first/publish/mms.service
[Unit]
Description=MiniMessanger ASP.NET WEB API (second)
After=multi-user.target
Requires=mysql.service

[Service]
WorkingDirectory=/var/www/mm
ExecStart=sudo /var/www/mm/second
Restart=always

[Install]
WantedBy=multi-user.target

